# build/base/etc/s6-overlay/s6-rc.d/vdr-prepare/up
#!/command/execlineb -P
define dir "/etc/s6-overlay/scripts"
define script "vdr-prepare_up"
exec "${dir}/${script}"
