# build/base/var/spool/cron/crontabs/vdr
0 4 * * * /usr/local/bin/naludumper-cron
