# build/base/etc/s6-overlay/s6-rc.d/vdr/run
#!/command/execlineb -P
with-contenv
foreground { s6-echo "┄┄┄┄┄┄┄┄┄┄┄┄┄┄┄┄┄┄┄┄┄┄┄┄┄┄⚟ Start VDR ⚞┄┄┄┄┄┄┄┄┄┄┄┄┄┄┄┄┄┄┄┄┄┄┄┄┄┄" }
emptyenv -p
export HOME /vdr
s6-setuidgid vdr
fdmove -c 2 1
/usr/bin/vdr
