# build/base/etc/s6-overlay/s6-rc.d/cronjobs/up
#!/command/execlineb -P
with-contenv
importas -D "true" start START_CROND
define dir "/etc/s6-overlay/scripts"
define script "cronjobs_up"
if -t { eltest $start != false }
 exec "${dir}/${script}"
