# build/base/etc/s6-overlay/s6-rc.d/init-user/up
#!/command/execlineb -P
with-contenv
define user vdr
define group vdr
backtick -E uid { id -u ${user} }
backtick -E gid {
 pipeline { getent group ${group} }
 s6-cut -d ":" -f 3
}
importas -D "666" puid PUID
importas -D "666" pgid PGID
foreground {
 if { eltest ${uid} -ne ${puid} }
  foreground { s6-echo "Change the user ID of user ${user} to ${puid}" }
  foreground { usermod --non-unique --uid "${puid}" "${user}" }
  if -n { redirfd -w 2 /dev/null find / -user "${uid}" -not -path "/proc/*" -exec chown -h "${user}" {} \; }
   s6-true
}
 if -t { eltest ${gid} -ne ${pgid} }
  foreground { s6-echo "Change the group ID of group ${group} to ${pgid}" }
  foreground { groupmod --non-unique --gid "${pgid}" "${group}" }
  if -nt { redirfd -w 2 /dev/null find / -group "${gid}" -not -path "/proc/*" -exec chgrp -h "${group}" {} \; }
   s6-true
