# build/base/etc/s6-overlay/s6-rc.d/updpkg/up
#!/command/execlineb -P
foreground { s6-echo "Download fresh package databases from the server..." }
/usr/sbin/pacman --quiet --noprogressbar -Sy
