# build/base/etc/s6-overlay/s6-rc.d/init-updpkg/up
#!/command/execlineb -P
foreground { s6-echo "Load current pacman package databases from the server" }
redirfd -w 1 /dev/null pacman --noprogressbar -Sy
