# build/base/etc/s6-overlay/s6-rc.d/init-updpkg/up
#!/command/execlineb -P
with-contenv
importas -D "false" no_inet NO_INET
if -nt { eltest ${no_inet} = true }
 foreground { s6-echo "Load current pacman package databases from the server" }
 if -nt { redirfd -w 1 /dev/null pacman --noprogressbar -Sy }
  s6-echo "WARNING: Update of the pacman package databases failed!"
