# root/etc/services.d/vdr/run
#!/usr/bin/with-contenv bash

echo "Starting vdr..."
export LANG="$LANG"
exec \
	s6-setuidgid vdr /usr/bin/vdr
