# build/base/etc/s6-overlay/s6-rc.d/crond/run
#!/command/execlineb -P
with-contenv
backtick -D "true" start { printcontenv START_CROND }
importas -usC start start
foreground {
 ifelse { test $start = false } {
  foreground { s6-echo crond is disabled by environment variable! }
  foreground { s6-rmrf /etc/s6-overlay/s6-rc.d/user/contents.d/crond }
  s6-svc -d -O /run/s6-rc/servicedirs/crond
 }
 if -n { test -d /var/spool/cron/crontabs }
  s6-mkdir -p /var/spool/cron/crontabs
}
if { test $start != false }
 /usr/sbin/busybox crond -f -l 5
