# build/base/etc/s6-overlay/s6-rc.d/crond/run
#!/command/execlineb -P
with-contenv
importas -D "true" start START_CROND
ifelse { eltest $start = false } {
  foreground { s6-echo "Cron-Daemon (crond) is disabled by environment variable!" }
  foreground { s6-rmrf /etc/s6-overlay/s6-rc.d/user/contents.d/crond }
  s6-svc -d -O /run/s6-rc/servicedirs/crond
 }
 foreground {
  if -n { eltest -d /var/spool/cron/crontabs }
   s6-mkdir -p /var/spool/cron/crontabs
 }
 /usr/sbin/busybox crond -f -l 5
