# build/base/etc/s6-overlay/s6-rc.d/init-permissions/up
#!/command/execlineb -P
define dir "/etc/s6-overlay/scripts"
define script "init-permissions_up"
exec "${dir}/${script}"
