# build/base/etc/s6-overlay/s6-rc.d/httpd/run
#!/command/execlineb -P
with-contenv
importas -D "false" start START_IMAGESERVER
ifelse { s6-test $start = false }
 {
  foreground { s6-echo "ImageServer (httpd) is not enabled via environment variable" }
  foreground { s6-rmrf /etc/s6-overlay/s6-rc.d/user/contents.d/httpd }
  s6-svc -d -O /run/s6-rc/servicedirs/httpd
 }
 s6-setuidgid http
 fdmove -c 2 1
 /usr/sbin/busybox httpd -f -p 8099
