# build/base/etc/s6-overlay/s6-rc.d/greeter/up
#!/command/execlineb -P
foreground { s6-echo "╔════════════════════════════════════════════════════════════════╗" }
foreground { s6-echo "║                        _  _  ____  ____                        ║" }
foreground { s6-echo "║                       ( \\/ )(  _ \\(  _ \\                       ║" }
foreground { s6-echo "║                        \\  /  )(_) ))   /                       ║" }
foreground { s6-echo "║                         \\/  (____/(_)\\_)                       ║" }
foreground { s6-echo "║                        http://www.tvdr.de                      ║" }
s6-echo "╚════════════════════════════════════════════════════════════════╝"
