# build/base/etc/s6-overlay/s6-rc.d/init/up
#!/command/execlineb -P
elglob -0s files /etc/s6-overlay/scripts/init.d/up_*
forx file { ${files} }
  importas file file
  ${file}
