# build/base/etc/s6-overlay/s6-rc.d/checkrec/up
#!/command/execlineb -P
if -t { s6-test -e /vdr/recordings/checkrec }
foreground { s6-echo "Starting CheckRecordings..." }
s6-setuidgid vdr
/usr/local/bin/checkrec -ql
