# build/base/etc/s6-overlay/s6-rc.d/checkrec/up
#!/command/execlineb -P
ifelse { test -e /vdr/recordings/checkrec }
 {
  foreground { s6-echo Starting CheckRecordings... }
  s6-setuidgid vdr
  fdmove -c 2 1
  /usr/local/bin/checkrec -ql
 }
 exit 0
